# supervisord.conf
[program:node]
command = /home/dotcloud/current/node_modules/.bin/coffee app.coffee
directory = /home/dotcloud/current
